# configs/gunicorn_multipage.socket
[Unit]
Description=gunicorn socket for multipage project

[Socket]
ListenStream=/run/gunicorn/multipage_socket

[Install]
WantedBy=sockets.target
